# yellowlight.conf
[program:yellowlight]
command=/usr/local/etc/yellowlight/yellowlight.py
numprocs=1
autostart=true
autorestart=true